# openstack-keystone.upstart
# Licensed under the Apache License, Version 2.0
# http://bazaar.launchpad.net/~ubuntu-branches/ubuntu/quantal/keystone/quantal/view/head:/debian/copyright
description "OpenStack Identity Service"
author "Soren Hansen <soren@linux2go.dk>"
# modified for RHEL6 by apevec@gmail.com

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec keystone-all" keystone
